# deploy/homelab-dash-backend.service
[Unit]
Description=Homelab Dashboard Backend
After=network.target

[Service]
Type=simple
User=homelab
Group=homelab
WorkingDirectory=/opt/homelab-dash/backend
Environment="PATH=/opt/homelab-dash/backend/.venv/bin:/usr/bin"
ExecStart=/opt/homelab-dash/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
